# .config/systemd/user/bg.service
[Unit]
Description=sets background

[Service]
ExecStart=/home/grzebiel/bin/bg-update

[Install]
WantedBy=default.target
