# .config/systemd/user/bg.timer
[Unit]
Description=refresh wallpaper

[Timer]
OnBootSec=30
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
