# .config/systemd/user/notesync.service
[Unit]
Description=sync notes with repo

[Service]
ExecStart=/home/grzebiel/bin/notesync

[Install]
WantedBy=default.target
