# .config/systemd/user/worktime.service
[Unit]
Description=Worktime monitoring

[Service]
ExecStart=/home/grzebiel/bin/worktime
Type=oneshot

[Install]
WantedBy=default.target
